# solarautopilot.service
[Unit]
Description=SolarAutopilot Website
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/solarautopilot
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
